# deploy/security-scanner.service
[Unit]
Description=Linux Security Scanner Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/security-scanner --config /etc/security-scanner/config.yml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
